# packages/bridge/systemd/install-update-from-usb.service.in
[Unit]
Description=Copy Update From USB Stick and Execute
After=syslog.target network.target chkphy.service systemd-modules-load.service
Before=bbbb.service

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/install-update.sh

[Install]
WantedBy=multi-user.target

