# packages/bridge/systemd/reset-playcontroller.service.in
[Unit]
Description=safety reset the playcontroller (for older motherboards), obtain and save the Unique Hardware ID
After=systemd-modules-load.service
Before=install-update-from-usb.service

[Service]
Type=oneshot
ExecStart=-@CMAKE_INSTALL_FULL_BINDIR@/reset-playcontroller.sh
ExecStart=-@CMAKE_INSTALL_FULL_BINDIR@/display-splashscreen.sh
