# package/audio-engine/system/fixCpuBindings.service.in
[Unit]
Description=Fix CPU bindings
After=syslog.target persistent.mount audio-engine.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/fixCpuBindings.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
